# src/main/java/cn/mysic/shell/cs/boot.conf
# /etc/init/sw/boot.conf

description "boot"

start on started mysql
#start on started sw/kafka

task

script
    cd /usr/share/server
    ./boot.sh
end script
